# packages/ublue-setup-services/src/services/ublue-system-setup.service
[Unit]
Description=Configure system
After=rpm-ostreed.service
Before=systemd-user-sessions.service

[Service]
Type=simple
ExecStart=/usr/libexec/ublue-system-setup
User=root

[Install]
WantedBy=multi-user.target
