# packages/ublue-os-media-automount-udev/ublue-os-media-automount.service
[Unit]
Description=Mount partitons automaticaly
After=local-fs.target
Before=graphical.target

[Service]
ExecStart=/usr/libexec/ublue-os-media-automount.py
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
