# packages/ublue-os-update-reminder/src/usr/lib/systemd/user/ublue-reboot-notifier.service
[Unit]
Description=Monitor screen unlock and show reboot notifications
After=graphical-session.target
ConditionPathExists=/run/user/%i/ublue-reboot-notifier-enabled

[Service]
Type=simple
ExecStart=/usr/libexec/ublue-reboot-notifier-watcher.py
Restart=on-failure
RestartSec=5
StartLimitBurst=3
StartLimitIntervalSec=600

[Install]
WantedBy=default.target
