# packages/ublue-os-update-reminder/src/usr/lib/systemd/user/ublue-uptime-checker-boot.service
[Unit]
Description=Disable reboot notifier on fresh boot

[Service]
Type=oneshot
ExecStart=/usr/libexec/ublue-uptime-checker.sh
RemainAfterExit=yes

[Install]
WantedBy=default.target
