# packages/ublue-os-update-reminder/src/usr/lib/systemd/user/ublue-uptime-checker.service
[Unit]
Description=Check system uptime and enable reboot notifier if needed

[Service]
Type=oneshot
ExecStart=/usr/libexec/ublue-uptime-checker.sh
