# packages/ublue-os-update-reminder/src/usr/lib/systemd/user/ublue-reboot-notifier.service
[Unit]
Description=Watch for unlock events
After=graphical-session.target

[Service]
ExecStart=/usr/libexec/ublue-reboot-notifier-watcher.py
Restart=always

[Install]
WantedBy=default.target
