# packages/ublue-setup-services/src/user-services/ublue-user-setup.service
[Unit]
Description=Configure system for current user
Wants=network-online.target
After=network-online.target
ConditionUser=!@system

[Service]
Type=simple
ExecStart=/usr/libexec/ublue-user-setup

[Install]
WantedBy=graphical-session.target
